# crontab
PATH=/usr/bin
BUSTIME_API_KEY=
PG_DATABASE=
# PG_HOST=""
# PYTHON=

*/2 * * * * make -e -C $HOME/mta-bus-archive positions
# */3 * * * * make -e -C $HOME/mta-bus-archive tripupdates
# */3 * * * * make -e -C $HOME/mta-bus-archive alerts
