# crontab
PATH=/usr/bin
DATABASE=
BUSTIME_API_KEY=

*/2 * * * * make -e -C mta-bus-archive positions
# */3 * * * * make -e -C mta-bus-archive tripupdates
# */3 * * * * make -e -C mta-bus-archive alerts
